# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=10
startsecs=1
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:next]
command=node /app/server.js
directory=/app
environment=NODE_ENV="production",HOSTNAME="127.0.0.1",PORT="3001"
autostart=true
autorestart=true
priority=20
startsecs=2
stopsignal=TERM
stopwaitsecs=20
user=nextjs
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:corescope-server]
command=/usr/local/bin/corescope-server -config-dir %(ENV_CORESCOPE_CONFIG_DIR)s -db %(ENV_CORESCOPE_DB_PATH)s -public %(ENV_CORESCOPE_PUBLIC_DIR)s -port %(ENV_CORESCOPE_PORT)s
directory=%(ENV_CORESCOPE_CONFIG_DIR)s
environment=DB_PATH="%(ENV_CORESCOPE_DB_PATH)s",CORESCOPE_INGESTOR_STATS="%(ENV_CORESCOPE_INGESTOR_STATS)s"
autostart=true
autorestart=true
priority=30
startsecs=2
stopsignal=TERM
stopwaitsecs=20
user=nextjs
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:corescope-ingestor]
command=/usr/local/bin/corescope-ingestor -config %(ENV_CORESCOPE_CONFIG_DIR)s/config.json
directory=%(ENV_CORESCOPE_CONFIG_DIR)s
environment=DB_PATH="%(ENV_CORESCOPE_DB_PATH)s",CORESCOPE_INGESTOR_STATS="%(ENV_CORESCOPE_INGESTOR_STATS)s"
autostart=false
autorestart=true
priority=40
startsecs=2
startretries=10
stopsignal=TERM
stopwaitsecs=20
user=nextjs
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
